# data/vm/cookbooks/supervisor/templates/default/grunt_watch.conf.erb
[program:grunt_watch]
command=bash -c "cd /vagrant && /usr/bin/grunt watch"
autostart=true
autorestart=true
stderr_logfile=/var/log/grunt_watch.err.log
stdout_logfile=/var/log/grunt_watch.out.log