# data/vm/cookbooks/supervisor/templates/default/grunt_watch.conf.erb
[program:grunt_watch]
user=vagrant
directory=/vagrant
command=/usr/bin/grunt watch
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/grunt_watch.err.log
stdout_logfile=/var/log/grunt_watch.out.log

; If there are any problems have a look at
;   * http://stackoverflow.com/q/9034709